# myscript.timer
[Unit]
Description=Timer for My Script Service

[Timer]
OnUnitActiveSec=5m
Unit=myscript.service

[Install]
WantedBy=timers.target
